# services/systemd/sensor_state_sync_consumer.service
[Unit]
Description=MoniCrowd Sensor State Consumer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/etc/telegraf/scripts
ExecStart=/usr/bin/python3 /etc/telegraf/scripts/sensor_state_sync_consumer.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
